# contrib/init/noirsharesd.conf
description "NoirShares Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env NOIRSHARESD_BIN="/usr/bin/noirsharesd"
env NOIRSHARESD_USER="noirshares"
env NOIRSHARESD_GROUP="noirshares"
env NOIRSHARESD_PIDDIR="/var/run/noirsharesd"
# upstart can't handle variables constructed with other variables
env NOIRSHARESD_PIDFILE="/var/run/noirsharesd/noirsharesd.pid"
env NOIRSHARESD_CONFIGFILE="/etc/noirshares/noirshares.conf"
env NOIRSHARESD_DATADIR="/var/lib/noirsharesd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # noirsharesd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$NOIRSHARESD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run noirsharesd."
        echo "The setting must appear in $NOIRSHARESD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"NoirShares Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$NOIRSHARESD_PIDDIR"
    chmod 0755 "$NOIRSHARESD_PIDDIR"
    chown $NOIRSHARESD_USER:$NOIRSHARESD_GROUP "$NOIRSHARESD_PIDDIR"
    chown $NOIRSHARESD_USER:$NOIRSHARESD_GROUP "$NOIRSHARESD_CONFIGFILE"
    chmod 0660 "$NOIRSHARESD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$NOIRSHARESD_PIDFILE" \
    --chuid $NOIRSHARESD_USER:$NOIRSHARESD_GROUP \
    --exec "$NOIRSHARESD_BIN" \
    -- \
    -pid="$NOIRSHARESD_PIDFILE" \
    -conf="$NOIRSHARESD_CONFIGFILE" \
    -datadir="$NOIRSHARESD_DATADIR" \
    -disablewallet \
    -daemon

